# dot_config/hypr/services/hourly-notify.service
[Unit]
Description=Notify every full hour

[Service]
Type=oneshot
ExecStart=/home/rasmus/.config/hypr/scripts/hourly-notify.sh
