# dot_config/hypr/timers/hourly-notify.timer
[Unit]
Description=Run hourly-notify every hour on the hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
